# supervisord.conf
[supervisord]
nodaemon=true

[program:paper_server]
directory=/etc/server
command=java -jar paper-server.jar
autorestart=true